# claude/launchagents/com.godl1ke.claude.mac-cleanup-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.godl1ke.claude.mac-cleanup-scan</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Use Homebrew bash 4+: mac-cleanup-scan.sh uses `declare -A` which
             crashes under Apple's frozen bash 3.2 (`declare: -A: invalid option`,
             exit 512). The /opt/homebrew/opt/bash/bin/bash path is keg-stable
             and survives `brew upgrade bash` symlink updates. -->
        <string>/opt/homebrew/opt/bash/bin/bash</string>
        <string>/Users/godl1ke/.dotfiles/claude/crons/mac-cleanup-scan.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/godl1ke</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/godl1ke/Library/Logs/claude-crons/mac-cleanup-scan-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/godl1ke/Library/Logs/claude-crons/mac-cleanup-scan-launchd-err.log</string>
</dict>
</plist>
