# claude/launchagents/com.godl1ke.claude.healthcheck-postrun.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.godl1ke.claude.healthcheck-postrun</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/godl1ke/.dotfiles/claude/crons/healthcheck.sh</string>
        <string>postrun</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/godl1ke</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/godl1ke/Library/Logs/claude-crons/healthcheck-postrun-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/godl1ke/Library/Logs/claude-crons/healthcheck-postrun-launchd-err.log</string>
</dict>
</plist>
